# init.d/chat-service
#!/sbin/runscript

PYTHON=python3
CONFIGFILE=/etc/chat_service.conf
PIDFILE=/var/run/chat-service.pid
USER=nobody

depend() {
  use messagehub-service
}

start() {
  ebegin "Starting ${SVCNAME}"
  start-stop-daemon --start --exec ${PYTHON} -b -m --name ${PYTHON} \
     --pidfile ${PIDFILE} -u ${USER} -- -m chat.chat_service -c ${CONFIGFILE}
  eend $?
}

stop() {
  ebegin "Stopping ${SVCNAME}"
  pkill -TERM -P `cat ${PIDFILE}`
  start-stop-daemon --stop --exec ${PYTHON} --name ${PYTHON} \
     --pidfile ${PIDFILE}
  eend $?
}
